# matugen-llava-seed/systemd/matugen-llava-seed.service
[Unit]
Description=matugen-llava-seed apply

[Service]
Type=oneshot
ExecStart=%h/.local/bin/matugen-llava-seed apply --image %f
