# matugen-llava-seed/systemd/matugen-llava-seed.path
[Unit]
Description=Watch wallpapers for changes

[Path]
PathModified=%h/Pictures/wallpapers
Unit=matugen-llava-seed.service

[Install]
WantedBy=default.target
